# deployment/integraldb.service
# Deployment configuration for backend service

[Unit]
Description=IntegralDB Backend Service
After=network.target

[Service]
Type=simple
User=integraldb
Group=integraldb
WorkingDirectory=/opt/integraldb
Environment="PATH=/opt/integraldb/venv/bin"
ExecStart=/opt/integraldb/venv/bin/integraldb-backend
Restart=always

[Install]
WantedBy=multi-user.target